# deploy/systemd/tock-server.service
[Unit]
Description=tock sync server (AGPL-3.0)
Documentation=https://github.com/kafkade/tock
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tock
Group=tock
ExecStart=/usr/local/bin/tock-server --bind 0.0.0.0:8080 --data-dir /var/lib/tock-server
Restart=on-failure
RestartSec=5

# Security hardening.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/tock-server
PrivateTmp=yes

# Environment.
Environment=RUST_LOG=info
# Account / auth configuration (see docs/self-hosting.md):
Environment=TOCK_MODE=self-hosted
Environment=TOCK_REGISTRATION_POLICY=disabled
# Optional: bootstrap an admin invite for this username on a fresh instance
# (the setup token is written to the journal). Leave unset to use the web
# first-run wizard instead.
#Environment=TOCK_ADMIN_USERNAME=admin@example.com

[Install]
WantedBy=multi-user.target
